# meta-python_st/.history/recipes-python_st/imx6sx-ksz/startksz8463/ksz8463_20210422143957.service
[Unit]
Description=KEPM service to start A9 application
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/python /ksz/sel_copp_fib.py
StandardOutput=tty

[Install]
WantedBy=sysinit.target