# meta-python_st/.history/recipes-python_st/imx6sx-ksz/startksz8463/ksz8463_20210422144116.service
[Unit]
Description=KEPM service to start A9 application
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
Type=simple
ExecStart=/usr/bin/python /ksz/sel_copp_fib.py
StandardInput=tty-force

[Install]
WantedBy=multi-user.target